# neuronos-iso/airootfs/etc/skel/.config/autostart/neuronos-onboarding.desktop
[Desktop Entry]
Type=Application
Name=NeuronOS Onboarding
Comment=Post-install onboarding wizard for NeuronOS
Exec=/usr/bin/neuronos-onboarding
Icon=neuronos
Terminal=false
Categories=System;Settings;
StartupNotify=true
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Phase=Application
